# rebuild-split-gateway.path
[Path]
PathChanged=/etc/ufw/user.rules
Unit=rebuild-split-gateway.service
