# watcher-split-gateway.path
[Path]
PathModified=/etc/ufw/user.rules

[Install]
WantedBy=multi-user.target